# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:api-gateway]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8001
directory=/app/services/api-gateway
autostart=true
autorestart=true
stderr_logfile=/var/log/api-gateway.err.log
stdout_logfile=/var/log/api-gateway.out.log
environment=PORT="8001"

[program:rag-engine]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8002
directory=/app/services/rag-engine
autostart=true
autorestart=true
stderr_logfile=/var/log/rag-engine.err.log
stdout_logfile=/var/log/rag-engine.out.log
environment=PORT="8002"

[program:visualization]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8003
directory=/app/services/visualization
autostart=true
autorestart=true
stderr_logfile=/var/log/visualization.err.log
stdout_logfile=/var/log/visualization.out.log
environment=PORT="8003"

[program:frontend]
command=npm start
directory=/app/frontend
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend.err.log
stdout_logfile=/var/log/frontend.out.log
environment=PORT="3000" 